# utrace/service/utrace.service
[Unit]
Description=Trace the uos exe.
After=network.target

[Service]
ExecStart=/opt/utrace/utrace

[Install]
WantedBy=default.target
