# utrace/service/utrace.service
[Unit]
Description=Trace the uos exe.
After=network.target

[Service]
ExecStart=/opt/myscript/utrace.sh

[Install]
WantedBy=default.target
